# automation/borg/systemd/borg-runner.service
[Unit]
Description=Multi-repo Borg Backup Engine
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/borg-runner.sh
